# claude-remote.service
[Unit]
Description=Claude Remote Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/jamie/projects/remote-claude-real
ExecStart=/home/jamie/.nix-profile/bin/pnpm tsx server.ts
Environment=NODE_ENV=production
Restart=always
RestartSec=3

# Logging to file
StandardOutput=append:/home/jamie/projects/remote-claude-real/logs/daemon-server.log
StandardError=append:/home/jamie/projects/remote-claude-real/logs/daemon-server.log

[Install]
WantedBy=default.target
