# claude-remote.service
[Unit]
Description=Claude Remote Server
After=network.target

[Service]
Type=simple
# Update these paths to match your installation
WorkingDirectory=%h/claude-remote
ExecStart=pnpm tsx server.ts
Environment=NODE_ENV=production
Restart=always
RestartSec=3

# Logging
StandardOutput=append:%h/claude-remote/logs/daemon-server.log
StandardError=append:%h/claude-remote/logs/daemon-server.log

[Install]
WantedBy=default.target
